# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn --log-level debug app:app -b localhost:5000
directory=/usr/src/app